# ethcloud/provision/roles/journald-cloudwatch/templates/etc/profile.d/go-path.sh.j2
#!/usr/bin/env bash

export GOPATH="{{go_path}}"
export GOROOT="{{go_root}}"
export PATH=$GOPATH/bin:$PATH
